# packaging/systemd/bootcontrold.socket
[Unit]
Description=BootControl D-Bus Daemon Socket

[Socket]
ListenStream=/run/bootcontrold.sock

[Install]
WantedBy=sockets.target
